# etc/supervisor/hipache.conf
[program:hipache]
command=/home/vagrant/jacc/node_modules/hipache/bin/hipache -c /home/vagrant/jacc/etc/hipache-config.json
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
